# etc/systemd/system/fstrim.service.d/override.conf
[Unit]
 Description=Discard unused blocks on filesystems from /proc/self/mountinfo

 [Service]
 ExecStart=/usr/bin/fstrim --listed-in /proc/self/mountinfo --verbose --quiet-unsupported
