# system/etc/systemd/system/getty@tty1.service.d/arch-mage.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- arch-mage' --noclear --skip-login - $TERM
