# etc/systemd/system/paccache.service.d/override.conf
[Service]
ExecStart=/usr/bin/paccache -rk0
